# allinone/etc/supervisor/conf.d/solr.conf
[program:solr]
command=ls -la
priority=13
autorestart=false
startretries=0